# rootfs/etc/inittab
# /etc/inittab: Init configuration for SliTaz ARM
#
::sysinit:/etc/init.d/rcS

# SHell invocations for selected ttys
tty1::respawn:/sbin/getty 38400 tty1
tty2::respawn:/sbin/getty 38400 tty2
tty2::respawn:/sbin/getty 38400 tty3

# Stuff to do before rebooting or shutting down
::restart:/etc/init.d/shutdown.sh
::restart:/sbin/init
::shutdown:/etc/init.d/shutdown.sh
::ctrlaltdel:/bin/umount -a -r
::ctrlaltdel:/sbin/reboot
