# rootfs/etc/inittab
# /etc/inittab: init configuration for SliTaz ARM.
#
::sysinit:/etc/init.d/rcS
::respawn:-/bin/sh -l
#tty1::respawn:/sbin/getty 38400 tty1
tty2::respawn:/sbin/getty 38400 tty2
tty3::askfirst:-/bin/sh
::ctrlaltdel:/bin/umount -a -r
::ctrlaltdel:/sbin/reboot
